# trufonomics-models/scripts/com.thales.dailyforecast.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
                       "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thales.dailyforecast</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kluless/kairos/trufonomics-models/scripts/run_daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kluless/kairos/trufonomics-models</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/kluless</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kluless/Library/Logs/thales/forecast.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kluless/Library/Logs/thales/forecast.err.log</string>
</dict>
</plist>
